# linux/onedrive/onedrive-bisync.timer
[Unit]
Description=Run onedrive-bisync.service every 15 minutes

[Timer]
OnCalendar=*:0/15
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target
