# linux/onedrive/onedrive-rclone.service
[Unit]
Description=OneDrive mount with rclone
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# ...existing code...
ExecStart=/usr/bin/rclone mount onedrive:/ %h/OneDrive-mount \
  --vfs-cache-mode full \
  --vfs-cache-max-age 24h \
  --dir-cache-time 24h \
  --buffer-size 128M
# Run a bisync before stopping/unmounting the mount (uses flock lock)
ExecStopPre=/usr/bin/flock -n /run/user/%U/onedrive-bisync.lock /usr/bin/rclone bisync onedrive:/ %h/OneDrive.local --cache-dir=%h/.cache/rclone --log-file=%h/.config/rclone/onedrive-bisync.log
# allow the bisync time to finish during shutdown
TimeoutStopSec=6min
ExecStop=/bin/fusermount -u %h/OneDrive-mount
Restart=no

[Install]
WantedBy=default.target
# How to run this service
# 0. Ensure User=<value> is correct for the system you're using
# 1. Save file to ~/.config/systemd/user/onedrive-rclone.service
# 2. systemctl --user daemon-reload
# 3. systemctl --user enable onedrive-rclone
# 4. systemctl --user start onedrive-rclone

# Check if the system is running systemctl --user status onedrive-rclone

