# linux/onedrive/onedrive-bisync.service
[Unit]
Description=Two-way sync for OneDrive with rclone bisync
After=network-online.target onedrive-rclone.service
Wants=network-online.target
Requires=onedrive-rclone.service

[Service]
Type=oneshot
# Check if OneDrive is mounted before proceeding
ExecStartPre=/bin/bash -c 'if [ ! -d %h/OneDrive ]; then echo "OneDrive not mounted" && exit 1; fi'
# Initialize bisync on first run (required for rclone bisync)
ExecStartPre=/bin/bash -c 'if [ ! -f %h/.config/rclone/onedrive-bisync-initialized ]; then /usr/bin/rclone bisync onedrive:/ %h/OneDrive --resync && touch %h/.config/rclone/onedrive-bisync-initialized; fi'
# Run the actual bisync
ExecStart=/usr/bin/rclone bisync onedrive:/ %h/OneDrive --progress --log-file=%h/.config/rclone/onedrive-bisync.log
# Mark completion
ExecStartPost=/usr/bin/touch %h/.config/rclone/onedrive-last-sync
RemainAfterExit=yes

[Install]
WantedBy=default.target
