# linux/ssh-agent.service
[Unit]
Description=SSH Agent (user)

[Service]
Type=forking
# write socket path and pid to the runtime env file so shells can source it
ExecStart=/usr/bin/ssh-agent -a %t/ssh-agent.socket > %t/ssh-agent.env
ExecStartPost=/usr/bin/sleep 0.2
# try to add your key on start (don't fail the service if adding fails)
ExecStartPost=/usr/bin/ssh-add /home/aidan/.ssh/id_ed25519_github || true
ExecStop=/usr/bin/ssh-agent -k
ExecStopPost=/bin/rm -f %t/ssh-agent.env
RemainAfterExit=yes

[Install]
WantedBy=default.target